# .test/nvim/systemd/smoke.service
[Unit]
Description=Smoke

[Service]
ExecStart=/bin/true

[Install]
WantedBy=multi-user.target
